# relay-reader/relay-read-file.service
[Unit]
Description=Theia relay read file

[Service]
ExecStart=/usr/loca/bin/relay-read-file

[Install]
WantedBy=multi-user.target
